# addons/thread-observability/rootfs/etc/s6-overlay/s6-rc.d/core/run
#!/command/with-contenv bash
exec python3 -m thread_observability.main --service core
