# addons/thread-observability/rootfs/etc/services.d/core/run
#!/bin/bash
set -e

exec python3 -m thread_observability.main --service core
